# systemd/battery-zen.service
[Unit]
Description=Battery Zen Daemon
After=multi-user.target
Wants=network.target

[Service]
Type=simple
User=%i
ExecStart=/home/prajwal/github-repos/battery-zen/battery-zen run
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Create necessary directories
ExecStartPre=/bin/mkdir -p /home/prajwal/.local/state/battery-zen

[Install]
WantedBy=multi-user.target
